# config/gunicorn.conf
/etc/systemd/system/gunicorn.servicedescription "Gunicorn application server handling sirius_server"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
setuid www-data
setgid www-data
chdir /var/www/sirius-server

exec venv/bin/gunicorn --workers 3 --bind unix:/var/www/sirius-server/sirius.sock sirius_server_app.wsgi:application